# tankdata_summary_v0_2_4/rootfs/etc/services.d/test/run
#!/usr/bin/with-contenv bash
echo "▶️ Starte Webserver..."
exec python3 /app/src/main.py
