# tankdata_summary_v0_2_2/rootfs/etc/services.d/test/run
#!/usr/bin/with-contenv bash
echo "▶️ Starte Webserver..."
python3 /app/src/main.py
